# ansible_collections/laroueverte/groundzero/roles/reboot_protection/files/etc/profile.d/reboot_protection_env.sh
alias reboot='/home/root/scripts/reboot_protection.sh'
alias reboot_test='/home/root/scripts/reboot_protection.sh --test'