# ansible_collections/laroueverte/groundzero/roles/common/files/iptables/iptables-v2.service
[Unit]
Description=IPv4 firewall with iptables

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/sbin/systemd-iptables start

[Install]
WantedBy=basic.target