# ansible_collections/laroueverte/groundzero/roles/docker/files/etc/systemd/system/docker.service.d/iptables.dependency.conf
[Unit]
After=iptables.service
BindsTo=iptables.service
