# ansible_collections/laroueverte/groundzero/roles/addokserver/files/lib/systemd/system/addok.service
[Unit]
Description=Addok Service
 
[Service]
ExecStart=/etc/addok/venv/bin/addok serve
StandardOutput=null
 
[Install]
WantedBy=multi-user.target
Alias=addok.service